# systems/milk_v_duo/extra-files/etc/init.d/blink
#!/sbin/openrc-run

LED_GPIO=/sys/class/gpio/gpio440

start() {
  touch /tmp/blinking
  /opt/milkv/blink.sh &
}

stop() {
  rm -f /tmp/blinking
}
