# systems/milk_v_duo/extra-files/etc/init.d/usb-host
#!/sbin/openrc-run

start() {
	echo host > /proc/cviusb/otg_role
}

