# files/extra-files-noble/etc/systemd/system/systemd-networkd-wait-online.service.d/override.conf
[Service]
# see: https://github.com/systemd/systemd/issues/30345#issuecomment-1843091779
ExecStart=/lib/systemd/systemd-networkd-wait-online --any
